# support/lsb_init_script_evercam_vpnclient
### BEGIN INIT INFO
# Provides:          evercam-vpnclient
# Required-Start:    $network
# Required-Stop:     $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Evercam VPN Client
# Description:       Starts the Sofether VPN Client
### END INIT INFO
DAEMON=/opt/vpnclient/vpnclient
test -x $DAEMON || exit 0

case "$1" in
	start)
		$DAEMON start
		;;
	stop)
		$DAEMON stop
		;;
	restart)
		$DAEMON stop
		sleep 3
		$DAEMON start
                ;;
	*)
		echo "Usage: $0 {start|stop|restart}"
		exit 1
esac
exit 0
